# systemd/security-manager-slave.target
[Unit]
Description=security-manager slave sockets
DefaultDependencies=true

