# systemd/security-manager.service.in
[Unit]
Description=Start the security manager
ConditionVirtualization=!lxc
ConditionPathExists=!/var/.zone/enabled

[Service]
Type=notify
ExecStart=@BIN_INSTALL_DIR@/security-manager

Sockets=security-manager.socket

[Install]
WantedBy=multi-user.target
