# systemd/security-manager-master.target
[Unit]
Description=security-manager master sockets
DefaultDependencies=true

