# systemd/security-manager.service.in
[Unit]
DefaultDependencies=false
After=local-fs.target
Description=Start the security manager

[Service]
Type=notify
ExecStart=@BIN_INSTALL_DIR@/security-manager
Sockets=security-manager.socket
Restart=always

[Install]
WantedBy=sysinit.target
