# systemd/security-manager-cleanup.service.in
[Unit]
Description=Cleanup private sharing
Before=starter.service

[Service]
Type=oneshot
ExecStart = @BIN_INSTALL_DIR@/security-manager-cleanup

